# api/docker/php/supervisor.d/messenger-worker.ini
[program:messenger-consume]
command=/usr/local/bin/php /srv/api/bin/console messenger:consume -vv --limit=10 async
user=www-data
numprocs=1
startsecs=0
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d


[program:messenger-consume-postprocessing]
command=/usr/local/bin/php /srv/api/bin/console messenger:consume -vv --limit=10 async_post_processing
user=www-data
numprocs=1
startsecs=0
autostart=true
autorestart=true
process_name=%(program_name)s_%(process_num)02d